# net/samba424/pkg-plist
bin/cifsdd
bin/dbwrap_tool
bin/dumpmscat
bin/gentest
bin/ldbadd
bin/ldbdel
bin/ldbedit
bin/ldbmodify
bin/ldbrename
bin/ldbsearch
bin/locktest
bin/masktest
bin/mdsearch
bin/mvxattr
bin/ndrdump
bin/net
bin/nmblookup
bin/ntlm_auth
bin/oLschema2ldif
bin/pdbedit
bin/profiles
bin/regdiff
bin/regpatch
bin/regshell
bin/regtree
bin/rpcclient
bin/samba-log-parser
bin/samba-regedit
%%ADS%%bin/samba-tool
bin/sharesec
bin/smbcacls
bin/smbclient
bin/smbcontrol
bin/smbcquotas
bin/smbget
bin/smbpasswd
bin/smbspool
bin/smbstatus
bin/smbtar
bin/smbtree
bin/tdbbackup
bin/tdbdump
bin/tdbrestore
bin/tdbtool
bin/testparm
bin/vfstest
bin/wbinfo
bin/wspsearch
etc/rc.d/samba_server
include/samba4/charset.h
include/samba4/core/doserr.h
include/samba4/core/error.h
include/samba4/core/hresult.h
include/samba4/core/ntstatus.h
include/samba4/core/ntstatus_gen.h
include/samba4/core/werror.h
include/samba4/core/werror_gen.h
%%LDAP%%include/samba4/smb_ldap.h
%%LDAP%%include/samba4/smbldap.h
include/samba4/dcesrv_core.h
include/samba4/domain_credentials.h
include/samba4/gen_ndr/atsvc.h
include/samba4/gen_ndr/auth.h
include/samba4/gen_ndr/claims.h
include/samba4/gen_ndr/dcerpc.h
include/samba4/gen_ndr/drsblobs.h
include/samba4/gen_ndr/drsuapi.h
include/samba4/gen_ndr/krb5pac.h
include/samba4/gen_ndr/lsa.h
include/samba4/gen_ndr/misc.h
include/samba4/gen_ndr/nbt.h
include/samba4/gen_ndr/ndr_atsvc.h
include/samba4/gen_ndr/ndr_dcerpc.h
include/samba4/gen_ndr/ndr_drsblobs.h
include/samba4/gen_ndr/ndr_drsuapi.h
include/samba4/gen_ndr/ndr_krb5pac.h
include/samba4/gen_ndr/ndr_misc.h
include/samba4/gen_ndr/ndr_nbt.h
include/samba4/gen_ndr/ndr_samr.h
include/samba4/gen_ndr/ndr_svcctl.h
include/samba4/gen_ndr/ndr_svcctl_c.h
include/samba4/gen_ndr/netlogon.h
include/samba4/gen_ndr/samr.h
include/samba4/gen_ndr/security.h
include/samba4/gen_ndr/server_id.h
include/samba4/gen_ndr/svcctl.h
include/samba4/ldb_wrap.h
include/samba4/libsmbclient.h
include/samba4/lookup_sid.h
include/samba4/machine_sid.h
include/samba4/ndr.h
include/samba4/ndr/ndr_dcerpc.h
include/samba4/ndr/ndr_drsblobs.h
include/samba4/ndr/ndr_drsuapi.h
include/samba4/ndr/ndr_krb5pac.h
include/samba4/ndr/ndr_nbt.h
include/samba4/ndr/ndr_svcctl.h
include/samba4/netapi.h
include/samba4/passdb.h
include/samba4/rpc_common.h
include/samba4/samba/session.h
include/samba4/samba/version.h
include/samba4/share.h
include/samba4/smb2_lease_struct.h
include/samba4/smb3posix.h
include/samba4/smbconf.h
include/samba4/tdr.h
include/samba4/tsocket.h
include/samba4/tsocket_internal.h
include/samba4/util/attr.h
include/samba4/util/blocking.h
include/samba4/util/data_blob.h
include/samba4/util/debug.h
include/samba4/util/discard.h
include/samba4/util/fault.h
include/samba4/util/genrand.h
include/samba4/util/idtree.h
include/samba4/util/idtree_random.h
include/samba4/util/signal.h
include/samba4/util/substitute.h
include/samba4/util/tfork.h
include/samba4/util/time.h
include/samba4/util_ldb.h
include/samba4/wbclient.h
lib/nss_winbind.so.1
lib/nss_wins.so.1
lib/pam_winbind.so
%%CUPS%%libexec/samba/smbspool_krb5_wrapper
%%SAMBA4_LIBDIR%%/libdcerpc-binding.so
%%SAMBA4_LIBDIR%%/libdcerpc-binding.so.0
%%SAMBA4_LIBDIR%%/libdcerpc-server-core.so
%%SAMBA4_LIBDIR%%/libdcerpc-server-core.so.0
%%SAMBA4_LIBDIR%%/libndr-krb5pac.so
%%SAMBA4_LIBDIR%%/libndr-krb5pac.so.0
%%SAMBA4_LIBDIR%%/libndr-nbt.so
%%SAMBA4_LIBDIR%%/libndr-nbt.so.0
%%SAMBA4_LIBDIR%%/libndr-standard.so
%%SAMBA4_LIBDIR%%/libndr-standard.so.0
%%SAMBA4_LIBDIR%%/libndr.so
%%SAMBA4_LIBDIR%%/libndr.so.6
%%SAMBA4_LIBDIR%%/libnetapi.so
%%SAMBA4_LIBDIR%%/libnetapi.so.1
%%SAMBA4_LIBDIR%%/libsamba-errors.so
%%SAMBA4_LIBDIR%%/libsamba-errors.so.1
%%SAMBA4_LIBDIR%%/libsamba-passdb.so
%%SAMBA4_LIBDIR%%/libsamba-passdb.so.0
%%SAMBA4_LIBDIR%%/libsamba-util.so
%%SAMBA4_LIBDIR%%/libsamba-util.so.0
%%SAMBA4_LIBDIR%%/libsmbclient.so
%%SAMBA4_LIBDIR%%/libsmbclient.so.0
%%SAMBA4_LIBDIR%%/libsmbconf.so
%%SAMBA4_LIBDIR%%/libsmbconf.so.0
%%LDAP%%%%SAMBA4_LIBDIR%%/libsmbldap.so
%%LDAP%%%%SAMBA4_LIBDIR%%/libsmbldap.so.2
%%SAMBA4_LIBDIR%%/libwbclient.so
%%SAMBA4_LIBDIR%%/libwbclient.so.0
%%SAMBA4_LIBDIR%%/modules/idmap/autorid.so
%%SAMBA4_LIBDIR%%/modules/idmap/hash.so
%%SAMBA4_LIBDIR%%/modules/idmap/nss.so
%%SAMBA4_LIBDIR%%/modules/idmap/rid.so
%%SAMBA4_LIBDIR%%/modules/idmap/script.so
%%SAMBA4_LIBDIR%%/modules/idmap/tdb.so
%%SAMBA4_LIBDIR%%/modules/idmap/tdb2.so
%%SAMBA4_LIBDIR%%/modules/krb5/async_dns_krb5_locator.so
%%SAMBA4_LIBDIR%%/modules/krb5/winbind_krb5_locator.so
%%SAMBA4_LIBDIR%%/modules/ldb/asq.so
%%SAMBA4_LIBDIR%%/modules/ldb/ildap.so
%%SAMBA4_LIBDIR%%/modules/ldb/ldb.so
%%SAMBA4_LIBDIR%%/modules/ldb/ldbsamba_extensions.so
%%SAMBA4_LIBDIR%%/modules/ldb/paged_searches.so
%%SAMBA4_LIBDIR%%/modules/ldb/rdn_name.so
%%SAMBA4_LIBDIR%%/modules/ldb/sample.so
%%SAMBA4_LIBDIR%%/modules/ldb/server_sort.so
%%SAMBA4_LIBDIR%%/modules/ldb/skel.so
%%SAMBA4_LIBDIR%%/modules/ldb/tdb.so
%%SAMBA4_LIBDIR%%/modules/nss_info/hash.so
%%SAMBA4_LIBDIR%%/modules/vfs/acl_tdb.so
%%SAMBA4_LIBDIR%%/modules/vfs/acl_xattr.so
%%SAMBA4_LIBDIR%%/modules/vfs/aio_fork.so
%%SAMBA4_LIBDIR%%/modules/vfs/aio_pthread.so
%%SAMBA4_LIBDIR%%/modules/vfs/audit.so
%%SAMBA4_LIBDIR%%/modules/vfs/cacheprime.so
%%SAMBA4_LIBDIR%%/modules/vfs/cap.so
%%SAMBA4_LIBDIR%%/modules/vfs/catia.so
%%SAMBA4_LIBDIR%%/modules/vfs/commit.so
%%SAMBA4_LIBDIR%%/modules/vfs/crossrename.so
%%SAMBA4_LIBDIR%%/modules/vfs/default_quota.so
%%SAMBA4_LIBDIR%%/modules/vfs/dirsort.so
%%SAMBA4_LIBDIR%%/modules/vfs/expand_msdfs.so
%%SAMBA4_LIBDIR%%/modules/vfs/extd_audit.so
%%SAMBA4_LIBDIR%%/modules/vfs/fake_perms.so
%%SAMBA4_LIBDIR%%/modules/vfs/freebsd.so
%%SAMBA4_LIBDIR%%/modules/vfs/fruit.so
%%SAMBA4_LIBDIR%%/modules/vfs/full_audit.so
%%SAMBA4_LIBDIR%%/modules/vfs/linux_xfs_sgid.so
%%SAMBA4_LIBDIR%%/modules/vfs/media_harmony.so
%%SAMBA4_LIBDIR%%/modules/vfs/offline.so
%%SAMBA4_LIBDIR%%/modules/vfs/preopen.so
%%SAMBA4_LIBDIR%%/modules/vfs/readahead.so
%%SAMBA4_LIBDIR%%/modules/vfs/readonly.so
%%SAMBA4_LIBDIR%%/modules/vfs/recycle.so
%%SAMBA4_LIBDIR%%/modules/vfs/shadow_copy.so
%%SAMBA4_LIBDIR%%/modules/vfs/shadow_copy2.so
%%SAMBA4_LIBDIR%%/modules/vfs/shell_snap.so
%%SAMBA4_LIBDIR%%/modules/vfs/streams_depot.so
%%SAMBA4_LIBDIR%%/modules/vfs/streams_xattr.so
%%SAMBA4_LIBDIR%%/modules/vfs/syncops.so
%%SAMBA4_LIBDIR%%/modules/vfs/time_audit.so
%%SAMBA4_LIBDIR%%/modules/vfs/unityed_media.so
%%SAMBA4_LIBDIR%%/modules/vfs/virusfilter.so
%%SAMBA4_LIBDIR%%/modules/vfs/widelinks.so
%%SAMBA4_LIBDIR%%/modules/vfs/worm.so
%%SAMBA4_LIBDIR%%/modules/vfs/xattr_tdb.so
%%SAMBA4_LIBDIR%%/modules/vfs/zfsacl.so
%%SAMBA4_LIBDIR%%/private/libCHARSET3-private-samba.so
%%SAMBA4_LIBDIR%%/private/libLIBWBCLIENT-OLD-private-samba.so
%%SAMBA4_LIBDIR%%/private/libMESSAGING-SEND-private-samba.so
%%SAMBA4_LIBDIR%%/private/libMESSAGING-private-samba.so
%%SAMBA4_LIBDIR%%/private/libREG-FULL-private-samba.so
%%SAMBA4_LIBDIR%%/private/libRPC-SERVER-LOOP-private-samba.so
%%SAMBA4_LIBDIR%%/private/libRPC-WORKER-private-samba.so
%%LDAP%%%%SAMBA4_LIBDIR%%/private/libsmbldaphelper-private-samba.so
%%SAMBA4_LIBDIR%%/private/libaddns-private-samba.so
%%SAMBA4_LIBDIR%%/private/libads-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libasn1-private-samba.so
%%SAMBA4_LIBDIR%%/private/libasn1-private-samba.so
%%SAMBA4_LIBDIR%%/private/libasn1util-private-samba.so
%%SAMBA4_LIBDIR%%/private/libauth-private-samba.so
%%SAMBA4_LIBDIR%%/private/libauth-unix-token-private-samba.so
%%SAMBA4_LIBDIR%%/private/libauth4-private-samba.so
%%SAMBA4_LIBDIR%%/private/libauthkrb5-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcli-cldap-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcli-ldap-common-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcli-ldap-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcli-nbt-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcli-smb-common-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcli-spoolss-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcliauth-private-samba.so
%%SAMBA4_LIBDIR%%/private/libclidns-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcluster-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcmdline-contexts-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcmdline-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libcom-err-private-samba.so
%%SAMBA4_BUNDLED_CMOCKA%%%%SAMBA4_LIBDIR%%/private/libcmocka-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcom-err-private-samba.so
%%SAMBA4_LIBDIR%%/private/libcommon-auth-private-samba.so
%%SAMBA4_LIBDIR%%/private/libdbwrap-private-samba.so
%%SAMBA4_LIBDIR%%/private/libdcerpc-pkt-auth-private-samba.so
%%SAMBA4_LIBDIR%%/private/libdcerpc-private-samba.so
%%SAMBA4_LIBDIR%%/private/libdcerpc-samba-private-samba.so
%%SAMBA4_LIBDIR%%/private/libdcerpc-samba4-private-samba.so
%%SAMBA4_LIBDIR%%/private/libdcerpc-samr-private-samba.so
%%SAMBA4_LIBDIR%%/private/libdnsserver-common-private-samba.so
%%SAMBA4_LIBDIR%%/private/libdsdb-module-private-samba.so
%%SAMBA4_LIBDIR%%/private/libevents-private-samba.so
%%SAMBA4_LIBDIR%%/private/libflag-mapping-private-samba.so
%%SAMBA4_LIBDIR%%/private/libgenrand-private-samba.so
%%SAMBA4_LIBDIR%%/private/libgensec-private-samba.so
%%SAMBA4_LIBDIR%%/private/libgpext-private-samba.so
%%SAMBA4_LIBDIR%%/private/libgpo-private-samba.so
%%SAMBA4_LIBDIR%%/private/libgse-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libgss-preauth-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libgssapi-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libhcrypto-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libhdb-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libheimbase-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libheimntlm-private-samba.so
%%SAMBA4_LIBDIR%%/private/libhttp-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libhx509-private-samba.so
%%SAMBA4_LIBDIR%%/private/libidmap-private-samba.so
%%SAMBA4_LIBDIR%%/private/libinterfaces-private-samba.so
%%SAMBA4_LIBDIR%%/private/libiov-buf-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libkdc-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libkrb5-private-samba.so
%%SAMBA4_LIBDIR%%/private/libkrb5samba-private-samba.so
%%SAMBA4_LIBDIR%%/private/libldb-cmdline-private-samba.so
%%SAMBA4_LIBDIR%%/private/libldb-key-value-private-samba.so
%%SAMBA4_LIBDIR%%/private/libldb-private-samba.so
%%SAMBA4_LIBDIR%%/private/libldb-tdb-err-map-private-samba.so
%%SAMBA4_LIBDIR%%/private/libldb-tdb-int-private-samba.so
%%SAMBA4_LIBDIR%%/private/libldbsamba-private-samba.so
%%SAMBA4_LIBDIR%%/private/liblibcli-lsa3-private-samba.so
%%SAMBA4_LIBDIR%%/private/liblibcli-netlogon3-private-samba.so
%%SAMBA4_LIBDIR%%/private/liblibsmb-private-samba.so
%%SAMBA4_LIBDIR%%/private/libmessages-dgm-private-samba.so
%%SAMBA4_LIBDIR%%/private/libmessages-util-private-samba.so
%%SAMBA4_LIBDIR%%/private/libmscat-private-samba.so
%%SAMBA4_LIBDIR%%/private/libmsghdr-private-samba.so
%%SAMBA4_LIBDIR%%/private/libmsrpc3-private-samba.so
%%SAMBA4_LIBDIR%%/private/libndr-samba-private-samba.so
%%SAMBA4_LIBDIR%%/private/libndr-samba4-private-samba.so
%%SAMBA4_LIBDIR%%/private/libnet-keytab-private-samba.so
%%SAMBA4_LIBDIR%%/private/libnetif-private-samba.so
%%SAMBA4_LIBDIR%%/private/libnpa-tstream-private-samba.so
%%SAMBA4_LIBDIR%%/private/libnss-info-private-samba.so
%%SAMBA4_LIBDIR%%/private/libposix-eadb-private-samba.so
%%SAMBA4_LIBDIR%%/private/libprinter-driver-private-samba.so
%%SAMBA4_LIBDIR%%/private/libprinting-migrate-private-samba.so
%%SAMBA4_LIBDIR%%/private/libregistry-private-samba.so
%%SAMBA4_LIBDIR%%/private/libreplace-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libroken-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamba-cluster-support-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamba-credentials-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamba-debug-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamba-hostconfig-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamba-modules-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamba-net-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamba-policy-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamba-security-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamba-security-trusts-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamba-sockets-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamba3-util-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamdb-common-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsamdb-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsecrets3-private-samba.so
%%SAMBA4_LIBDIR%%/private/libserver-id-db-private-samba.so
%%SAMBA4_LIBDIR%%/private/libserver-role-private-samba.so
%%SAMBA4_LIBDIR%%/private/libshares-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsmbclient-raw-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsmbd-base-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsmbd-shim-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsmbpasswdparser-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsocket-blocking-private-samba.so
%%SAMBA4_LIBDIR%%/private/libstable-sort-private-samba.so
%%SAMBA4_LIBDIR%%/private/libsys-rw-private-samba.so
%%SAMBA4_LIBDIR%%/private/libtalloc-private-samba.so
%%SAMBA4_LIBDIR%%/private/libtalloc-report-printf-private-samba.so
%%SAMBA4_LIBDIR%%/private/libtalloc-report-private-samba.so
%%SAMBA4_LIBDIR%%/private/libtdb-private-samba.so
%%SAMBA4_LIBDIR%%/private/libtdb-wrap-private-samba.so
%%SAMBA4_LIBDIR%%/private/libtevent-private-samba.so
%%SAMBA4_LIBDIR%%/private/libtevent-util-private-samba.so
%%SAMBA4_LIBDIR%%/private/libtime-basic-private-samba.so
%%SAMBA4_LIBDIR%%/private/libtorture-private-samba.so
%%SAMBA4_LIBDIR%%/private/libutil-crypt-private-samba.so
%%SAMBA4_LIBDIR%%/private/libutil-reg-private-samba.so
%%SAMBA4_LIBDIR%%/private/libutil-setid-private-samba.so
%%SAMBA4_LIBDIR%%/private/libutil-tdb-private-samba.so
%%GSSAPI_BUILTIN%%%%SAMBA4_LIBDIR%%/private/libwind-private-samba.so
%%SAMBA4_LIBDIR%%/private/libxattr-tdb-private-samba.so
libdata/ldconfig/samba424
%%PKGCONFIGDIR%%/ndr.pc
%%PKGCONFIGDIR%%/ndr_krb5pac.pc
%%PKGCONFIGDIR%%/ndr_nbt.pc
%%PKGCONFIGDIR%%/ndr_standard.pc
%%PKGCONFIGDIR%%/netapi.pc
%%PKGCONFIGDIR%%/samba-util.pc
%%PKGCONFIGDIR%%/smbclient.pc
%%PKGCONFIGDIR%%/wbclient.pc
libexec/samba/rpcd_classic
libexec/samba/rpcd_epmapper
libexec/samba/rpcd_fsrvp
libexec/samba/rpcd_lsad
libexec/samba/rpcd_mdssvc
libexec/samba/rpcd_spoolss
libexec/samba/rpcd_winreg
libexec/samba/samba-bgqd
libexec/samba/samba-dcerpcd
sbin/eventlogadm
sbin/nmbd
sbin/smbd
sbin/winbindd
@dir %%SAMBA4_LIBDIR%%/modules/auth
@dir %%SAMBA4_LIBDIR%%/modules/bind9
@dir %%SAMBA4_LIBDIR%%/modules/gensec
@dir %%SAMBA4_LIBDIR%%/modules/gpext
@dir %%SAMBA4_LIBDIR%%/modules/pdb
@dir %%SAMBA4_LIBDIR%%/modules/perfcount
@dir %%SAMBA4_LIBDIR%%/modules/process_model
@dir %%SAMBA4_LIBDIR%%/modules/service
@dir %%SAMBA4_BINDDNSDIR%%
@dir %%SAMBA4_PRIVATEDIR%%
@dir %%SAMBA4_LOCKDIR%%
@dir %%SAMBA4_LOGDIR%%
@dir %%SAMBA4_RUNDIR%%
